# .claude/systemd/claude-refresh.service
[Unit]
Description=Claude Code orphaned process refresh
Documentation=https://github.com/anthropics/claude-code

[Service]
Type=oneshot
# Use the script from the ProofChecker project
# Update this path if you install elsewhere
ExecStart=%h/Projects/ProofChecker/.claude/scripts/claude-refresh.sh --force
# Alternatively, if installed system-wide:
# ExecStart=/usr/local/bin/claude-refresh.sh --force

# Run with normal user permissions
User=%u

# Environment
Environment=HOME=%h
Environment=PATH=/usr/bin:/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-refresh

[Install]
WantedBy=default.target
